# package/luci-app-ktm-tools/root/etc/uci-defaults/ktm-tools
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ktm-tools[-1]
	add ucitrack ktm-tools
	set ucitrack.@ktm-tools[-1].init=ktm-tools
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
